# systemd/bobsled@.service
# systemd/bobsled@.service
[Unit]
Description=bobsled GitHub Actions runner slot %i
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=%h/image-digest.env
RuntimeDirectory=bobsled/%i
RuntimeDirectoryMode=0700
ExecStartPre=%h/.local/bin/bobsled-mint \
    --instance %i \
    --output %t/bobsled/%i/jit.json
ExecStart=/usr/bin/podman run --rm --replace \
    --name=bobsled-%i \
    --userns=keep-id:uid=1000,gid=1000 \
    --device=/dev/fuse \
    --tmpfs=/tmp \
    --add-host=host.containers.internal:host-gateway \
    -v %t/bobsled/%i:/jit:ro \
    -v %h/.cache/bobsled/slots/%i/current:/cache \
    -v %h/registries.conf:/home/podman/.config/containers/registries.conf:ro \
    bobsled:${BOBSLED_IMAGE_DIGEST}
ExecStopPost=/bin/rm -rf %t/bobsled/%i
Restart=always
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
